# system/archappl@engine.service.d/extra.conf
[Unit]
After=network.target archappl@mgmt.service archappl@etl.service
Requires=network.target archappl@mgmt.service archappl@etl.service
PartOf=archappl.service
