# system/archappl@etl.service.d/extra.conf
[Unit]
After=network.target archappl@mgmt.service
Requires=network.target archappl@mgmt.service
